# cron/crontab
# /etc/crontab に書き足す
*/5 * * * * root /home/isucon/isucon3/cron/crontab.sh 2>/tmp/cron.error.log
