# cron/crontab
*/5  *  *  *  * cd /home/isucon/isucon3 && /home/isucon/isucon3/cron/crontab.sh
